# systemd/halo-remote.service
[Unit]
Description=HALO Remote Server
After=network.target
StartLimitIntervalSec=0

[Service]
Type=exec
Restart=always
RestartSec=1
User=root
ExecStart=/usr/local/sbin/halo_remote
EnvironmentFile=-/etc/sysconfig/halo

[Install]
WantedBy=multi-user.target
